# contrib/systemd/gitforgefs@.service
[Unit]
Description=FUSE filesystem for gitlab groups and projects (%i)
Wants=network-online.target
After=network-online.target

[Service]
ExecStart=%h/go/bin/gitforgefs -config %E/gitforgefs/%i.yaml

[Install]
WantedBy=default.target